# user/editors.target
[Unit]
Description=Target for editors running in the background. (nvim, emacs etc.)

[Install]
WantedBy=default.target
Also=editor-services.target